# systemd/setu-midnight-jobs.service
[Unit]
Description=Setu V3 Midnight Data Collection
After=network.target postgresql.service

[Service]
Type=oneshot
User=pi50001_admin
WorkingDirectory=/home/pi50001_admin/SetuV3
Environment="PATH=/home/pi50001_admin/SetuV3/venv/bin"
ExecStart=/home/pi50001_admin/SetuV3/venv/bin/python tools/run_midnight_jobs.py
StandardOutput=journal
StandardError=journal
